# supervisord.conf
[supervisord]
nodaemon=true
logfile = /var/log/supervisord.log

[program:nginx]
command=nginx -g "daemon off;"

[program:swiftly]
command=gunicorn -t 600 -w 4 -b 0.0.0.0:5000 --paste production.ini

