# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g daemon off;

[program:swiftly]
command=gunicorn -t 600 -w 4 -u bradley -b 0.0.0.0:5000 --paste production.ini

